# app/Http/Controllers/Admin/config/ap1.ini
[program:wifi_fix]
command=/etc/init.d/bwifi start
autostart=1
autorestart=
priority=1
stderr_logfile=/var/log/wifi.out.log
[program:start_ap]
command=nmcli con up My-Hotspot
autostart=1
autorestart=
priority=2
stderr_logfile=/var/log/wifi.out.log
[program:reload_dhcpcd]
command=systemctl daemon-reload
autostart=1
autorestart=
priority=3
stderr_logfile=/var/log/wifi.out.log
